# puppet/files/etc/init/rq_scheduler.conf
# rq task scheduler worker

description "rq scheduler"
start on runlevel [2345]
stop on runlevel [06]

respawn
chdir /www/src
exec rqscheduler
