# puppet/files/etc/init/rqworker.conf
# rq  worker

description "rq worker"

setuid "www-data"
env HOME=/www/src

start on runlevel [2345]
stop on runlevel [06]

chdir /www/src
exec rqworker btc data
